# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/hotplug.d/button/00-reset-to-factory
#!/bin/sh

[ "${BUTTON}" != "reset" ] && exit 0
[ "${ACTION}" != "pressed" ] && exit 0

/usr/sbin/reset-to-factory
