# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/hotplug.d/button/09-service
#!/bin/sh

[ "${BUTTON}" != "wps" ] && exit 0
[ "${ACTION}" != "held" ] && exit 0

[ "$(nvram get web_auth)" = "passwd" ] && exit 0

status set web_autologin "on"

nbctl start_webautologin
