# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/hotplug.d/button/10-eco
#!/bin/sh

[ "${BUTTON}" != "BTN_1" ] && exit 0
[ "${ACTION}" != "pressed" ] && exit 0

/etc/init.d/eco switch
/etc/init.d/wifisched reload
