# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/hotplug.d/iface/05-qos
#!/bin/sh

[ "${INTERFACE}" != "data" ] && exit 0

case ${ACTION} in
	ifup)
		/etc/init.d/qos start
		;;
	ifdown)
		/etc/init.d/qos stop
		;;
esac
