# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/hotplug.d/iface/32-ipv6
#!/bin/sh

[ "${INTERFACE}" != "data" ] && exit 0

case ${ACTION} in
	ifup)
		/etc/init.d/ipv6 reload
		;;
	ifdown)
		;;
esac
