# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/669082237/rootfs_ubifs/etc/init.d/samba
#!/bin/sh /etc/rc.common

logger -t "samba[$$]" "$action"

make_cfg() {
	#
	# doc for samba configuration file:
	#
	# http://www.samba.org/samba/docs/man/Samba-HOWTO-Collection/speed.html
	#
	# - follow symlinks = no
	#
	#         don't follow symlinks. !! fix vulnerability issue !!
	#         Otherwise, you can access any folder on neufbox by creating 
	#         a symbolic link to that folder.
	#
	# - socket options = TCP_NODELAY
	#
	#        The socket option TCP_NODELAY is the one that seems to make 
	#        the biggest single difference for most networks. Many people 
	#        report that adding socket options = TCP_NODELAY doubles the 
	#        read performance of a Samba drive. The best explanation I 
	#        have seen for this is that the Microsoft TCP/IP stack is slow 
	#        in sending TCP ACKs.
	#

	echo "# - autogenerated -" > $1
	echo "[global]" >> $1
	[ "$(nvram get samba_debug)" = "on" ] && echo "debug level = 3" >> $1
	echo "workgroup = $(nvram get samba_workgroup)" >> $1
	echo "server string = $(status get operator_hostname) samba %v" >> $1
	echo "encrypt passwords = yes" >> $1
	echo "passdb backend = smbpasswd" >> $1
	echo "socket options = TCP_NODELAY" >> $1
	echo "netbios name = $(nvram get samba_name)" >> $1
	echo "load printers = no" >> $1
	echo "security = share" >> $1
	echo "bind interfaces only = True" >> $1
	echo "interfaces = $(status get lan_ifname)" >> $1
	echo "follow symlinks = no" >> $1

	local smb_config=/tmp/smb.xml
	local disk_config=/tmp/disk.xml
	nvram xml samba > ${smb_config}
	disk block_list > ${disk_config}
	local n=1
	for rule in $(roxml -q ${smb_config} '//samba/share/uuid')
	do
		eval $(roxml -o ${smb_config} "//samba/share[$n]")
		n=$(($n + 1))
		case "${name}" in
		"homes"|"printers"|"global"|"globals")
			logger -t "samba[$$]" "do not shared ${name}"
			;;
		*) 
			eval $(roxml -o ${disk_config} "//partition[uuid=\"${uuid}\"]")
			if [ ${mounted} = "1" ]; then
				echo "" >> $1
				echo "[${name}]" >> $1
				echo "comment = ${name}" >> $1
				echo "path = ${mntpath}${dir}" >> $1
				echo "browseable = yes" >> $1
				echo "writable = yes" >> $1
				echo "public = yes" >> $1
			else
				logger -t "samba[$$]" "disk uuid"\
					"${uuid} not mounted"
			fi
			unset mounted
			unset mntpath
			;;
		esac
	done
	rm -f ${smb_config}
	rm -f ${disk_config}
}

do_start() {
	[ "$(nvram get samba_active)" = "off" ] && return
	if [ "$(status get samba_status)" = "up" ];then
		logger -t "samba[$$]" "samba already running"
		return
	fi

	if [ $(hostname) = "nb4" ];then
		if [ "$(status get packages_samba)" != "up" ]; then
			lock -u /run/lock/samba.lock
			/etc/init.d/packages reload
			return
		fi
	fi
	
	status set samba_status starting

	mkdir -p /tmp/etc/samba
	if ! make_cfg /etc/samba/smb.conf; then
		status set samba_status error_unknown
		return
	fi

	smbd -D -s /etc/samba/smb.conf
	nmbd -D -s /etc/samba/smb.conf

	iptables -A LAN_FILTER -p udp --dport netbios-ns:netbios-ssn -j ACCEPT
	iptables -A LAN_FILTER -p tcp --dport netbios-ns:netbios-ssn -j ACCEPT
	iptables -A LAN_FILTER -p tcp --dport 445 -j REJECT

	status set samba_status up
}

do_stop() {
	iptables -D LAN_FILTER -p udp --dport netbios-ns:netbios-ssn -j ACCEPT
	iptables -D LAN_FILTER -p tcp --dport netbios-ns:netbios-ssn -j ACCEPT
	iptables -D LAN_FILTER -p tcp --dport 445 -j REJECT

	killall smbd
	killall nmbd
	sleep 1
	killall -KILL smbd 2>/dev/null
	killall -KILL nmbd 2>/dev/null

	status set samba_status down

	rm -f /etc/samba/smb.conf
}

start() {
	lock /run/lock/samba.lock
	do_start
	lock -u /run/lock/samba.lock
}

stop() {
	lock /run/lock/samba.lock
	do_stop
	lock -u /run/lock/samba.lock
}

restart() {
	lock /run/lock/samba.lock
	do_stop
	sleep 1
	do_start
	lock -u /run/lock/samba.lock
}
