# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/init.d/boot-terminated
#!/bin/sh /etc/rc.common

START=99

boot() {
	touch /var/state/services
	event notify boot-succeeded

if [ "$(nvram get eco_active)" = "on" ]; then
	/etc/init.d/eco start
fi
	ledctl power on
}
