# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/init.d/watchdog
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

	
START=11
start() {
        [ -c /dev/watchdog ] && [ -x /sbin/watchdog ] && \
        	watchdog -t 10 -T 85 /dev/watchdog
}
