# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/1953888591/rootfs_ubifs/etc/init.d/yaddns
#!/bin/sh /etc/rc.common

EXTRA_COMMANDS="forceupdate"

logger -t "yaddns[$$]" "$action"

make_cfg() {
	local ddns_service=$(nvram get ddns_service)
	local ddns_username=$(nvram get ddns_username)
	local ddns_password=$(nvram get ddns_password)
	local ddns_hostname=$(nvram get ddns_hostname)

	echo "# - autogenerated -" > $1
	echo "wanifname = \"$(status get net_data_ifname)\"" >> $1
	echo "mode = \"direct\"" >> $1
	echo "" >> $1
	echo "account {" >> $1
	echo "  name = \"${ddns_service} account\"" >> $1
	echo "  service = \"${ddns_service}\"" >> $1
	echo "  username = \"${ddns_username}\"" >> $1
	echo "  password = \"${ddns_password}\"" >> $1
	echo "  hostname = \"${ddns_hostname}\"" >> $1
	echo "}" >> $1
}

do_start() {
	if [ "$(nvram get ddns_active)" != "on" ]; then
		return 0
	fi

	if [ "$(status get net_data_infra)" = "gprs/ppp" ]; then
		logger -t "yaddns[$$]" "3g network. Stop it."
		status set ddns_status err_gprs
		return 0
	fi

	make_cfg /var/yaddns.cfg

	test -e /run/yaddns.pid && local pid=$(cat /run/yaddns.pid)
	if [ -n "${pid}" ]; then
		kill -HUP ${pid}
		return 0
	fi

	yaddns -D -L -f /var/yaddns.cfg -p /run/yaddns.pid
}

do_stop() {
	test -e /run/yaddns.pid && local pid=$(cat /run/yaddns.pid)
	if [ ! -n "${pid}" ]; then
		return 1
	fi

	kill -TERM ${pid}

	for i in $(seq 1 10); do
		sleep 1
		if ! kill -0 ${pid}; then
			rm -f /run/yaddns.pid
			return 0	
		fi
	done

	kill -KILL ${pid}
	rm -f /run/yaddns.pid
}

start() {
	lock /run/lock/yaddns.lock
	do_start
	lock -u /run/lock/yaddns.lock
}

stop() {
	lock /run/lock/yaddns.lock
	do_stop
	lock -u /run/lock/yaddns.lock
}

restart() {
	lock /run/lock/yaddns.lock
	do_stop
	sleep 1
	do_start
	lock -u /run/lock/yaddns.lock
}

reload() {
	lock /run/lock/yaddns.lock
	do_start
	lock -u /run/lock/yaddns.lock
}

forceupdate() {
	lock /run/lock/yaddns.lock
	test -e /run/yaddns.pid && local pid=$(cat /run/yaddns.pid)
	if [ -n "${pid}" ]; then
		make_cfg /var/yaddns.cfg
		kill -USR2 ${pid}
	else
		do_start
	fi
	lock -u /run/lock/yaddns.lock
}
