# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/669082237/rootfs_ubifs/etc/hotplug.d/iface/30-miniupnpd
#!/bin/sh

[ "${INTERFACE}" != "data" ] && exit 0

case ${ACTION} in
	ifup)
		/etc/init.d/miniupnpd restart
		;;
	ifdown)
		/etc/init.d/miniupnpd restart
		;;
esac
