# firmware/_cfe_bcm63xx_03022026.bin.extracted/ubifs-root/669082237/rootfs_ubifs/etc/hotplug.d/iface/31-yaddns
#!/bin/sh

[ "${INTERFACE}" != "data" ] && exit 0

case ${ACTION} in
	ifup)
		# try to reload yaddns. if doesn't started, it will be
		/etc/init.d/yaddns reload
		;;
	ifdown)
		;;
esac
